# config/includes.chroot/etc/xdg/autostart/noriva-installer.desktop
[Desktop Entry]
Type=Application
Name=Auto Installer
Exec=/usr/local/bin/noriva-autostart
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
